# deploy-configs/supervisor.conf
; ==============================================================================
# Laravel ERP Konstruksi - Supervisor Configuration (Complete)
; ==============================================================================
# Production-ready Supervisor config for queue workers, Horizon, and scheduler
; Alternative to systemd services - choose ONE approach (systemd OR supervisor)
; ==============================================================================
# Install: sudo cp deploy-configs/supervisor.conf /etc/supervisor/conf.d/erp-konstruksi.conf
# Reload:  sudo supervisorctl reread && sudo supervisorctl update
# Status:  sudo supervisorctl status
; ==============================================================================

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770
chown=root:www-data

[inet_http_server]
port=127.0.0.1:9001
username=admin
password=${SUPERVISOR_PASSWORD:-changeme}

[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=65535
minprocs=32768
user=root
directory=/tmp
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
username=admin
password=${SUPERVISOR_PASSWORD:-changeme}
prompt=supervisorctl>

; ==============================================================================
# LARAVEL HORIZON (Recommended for Redis Queues)
; ==============================================================================
# Horizon provides dashboard, metrics, and automatic scaling
# Enable if using Redis queues with Horizon package

[program:erp-horizon]
command=php /var/www/erp-konstruksi/artisan horizon
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
autorestart_retries=3
autorestart_delay=10
startsecs=10
startretries=3
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/erp-horizon.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
environment=APP_ENV=production,APP_DEBUG=false,PHP_MEMORY_LIMIT=512M
stopasgroup=true
killasgroup=true
stopwaitsecs=120

; ==============================================================================
# STANDARD QUEUE WORKERS (If NOT using Horizon)
; ==============================================================================
# Use these instead of Horizon for simpler setups or database queues

; High Priority Queue (critical jobs, notifications, webhooks)
[program:erp-queue-high]
command=php /var/www/erp-konstruksi/artisan queue:work redis --queue=high --sleep=3 --tries=3 --timeout=120 --max-jobs=1000 --memory=256 --rest=60
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
autorestart_retries=5
autorestart_delay=5
startsecs=10
startretries=3
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/erp-queue-high.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
environment=APP_ENV=production,APP_DEBUG=false
stopasgroup=true
killasgroup=true
stopwaitsecs=3600

; Default Queue (standard application jobs)
[program:erp-queue-default]
command=php /var/www/erp-konstruksi/artisan queue:work redis --queue=default --sleep=5 --tries=3 --timeout=180 --max-jobs=500 --memory=256 --rest=60
process_name=%(program_name)s_%(process_num)02d
numprocs=4
autostart=true
autorestart=true
autorestart_retries=5
autorestart_delay=5
startsecs=10
startretries=3
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/erp-queue-default.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
environment=APP_ENV=production,APP_DEBUG=false
stopasgroup=true
killasgroup=true
stopwaitsecs=3600

; Low Priority Queue (bulk operations, exports, reports)
[program:erp-queue-low]
command=php /var/www/erp-konstruksi/artisan queue:work redis --queue=low --sleep=10 --tries=2 --timeout=300 --max-jobs=200 --memory=128 --rest=120
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
autorestart_retries=3
autorestart_delay=10
startsecs=10
startretries=3
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/erp-queue-low.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
environment=APP_ENV=production,APP_DEBUG=false
stopasgroup=true
killasgroup=true
stopwaitsecs=3600

; RAB Import Queue (heavy Excel/CSV imports, long-running)
[program:erp-queue-rab-import]
command=php /var/www/erp-konstruksi/artisan queue:work redis --queue=rab-import --sleep=3 --tries=2 --timeout=600 --max-jobs=50 --memory=512 --rest=60
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
autorestart_retries=3
autorestart_delay=10
startsecs=10
startretries=3
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/erp-queue-rab-import.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
environment=APP_ENV=production,APP_DEBUG=false
stopasgroup=true
killasgroup=true
stopwaitsecs=7200

; AI Classification Queue (ML/AI processing jobs)
[program:erp-queue-ai-classify]
command=php /var/www/erp-konstruksi/artisan queue:work redis --queue=ai-classify --sleep=5 --tries=3 --timeout=1800 --max-jobs=100 --memory=512 --rest=120
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
autorestart_retries=3
autorestart_delay=15
startsecs=10
startretries=3
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/erp-queue-ai-classify.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
environment=APP_ENV=production,APP_DEBUG=false
stopasgroup=true
killasgroup=true
stopwaitsecs=7200

; Failed Jobs Handler (optional - processes failed jobs table)
[program:erp-queue-failed]
command=php /var/www/erp-konstruksi/artisan queue:work redis --queue=failed --sleep=60 --tries=1 --timeout=60 --max-jobs=10 --memory=128 --rest=300
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
autorestart_retries=3
autorestart_delay=60
startsecs=10
startretries=3
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/erp-queue-failed.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
environment=APP_ENV=production,APP_DEBUG=false
stopasgroup=true
killasgroup=true
stopwaitsecs=3600

; ==============================================================================
# LARAVEL SCHEDULER (Cron Replacement)
; ==============================================================================
# Runs every minute to execute scheduled tasks
; IMPORTANT: Only run ONE scheduler instance across all servers!

[program:erp-scheduler]
command=php /var/www/erp-konstruksi/artisan schedule:run --verbose --no-interaction
process_name=%(program_name)s
numprocs=1
autostart=true
autorestart=true
autorestart_retries=3
autorestart_delay=60
startsecs=5
startretries=3
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/erp-scheduler.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
environment=APP_ENV=production,APP_DEBUG=false
stopasgroup=true
killasgroup=true
stopwaitsecs=60

; Cron-style: Run every minute using cron instead of supervisor for scheduler
; Add to crontab (sudo crontab -e -u www-data):
; * * * * * cd /var/www/erp-konstruksi && php artisan schedule:run >> /dev/null 2>&1

; ==============================================================================
# WEBSOCKET SERVER (Laravel Reverb / Soketi / Pusher)
; ==============================================================================
; Enable if using real-time features (broadcasting, presence channels)

; [program:erp-reverb]
; command=php /var/www/erp-konstruksi/artisan reverb:start --host=0.0.0.0 --port=8080
; process_name=%(program_name)s
; numprocs=1
; autostart=true
; autorestart=true
; autorestart_retries=3
; autorestart_delay=10
; startsecs=5
; startretries=3
; user=www-data
; redirect_stderr=true
; stdout_logfile=/var/log/supervisor/erp-reverb.log
; stdout_logfile_maxbytes=20MB
; stdout_logfile_backups=5
; environment=APP_ENV=production,APP_DEBUG=false
; stopasgroup=true
; killasgroup=true
; stopwaitsecs=30

; ==============================================================================
# EVENT LISTENER WORKERS (For high-volume event processing)
; ==============================================================================

; [program:erp-event-listener]
; command=php /var/www/erp-konstruksi/artisan queue:listen redis --queue=events --sleep=3 --tries=3 --timeout=60 --memory=256
; process_name=%(program_name)s_%(process_num)02d
; numprocs=2
; autostart=true
; autorestart=true
; user=www-data
; redirect_stderr=true
; stdout_logfile=/var/log/supervisor/erp-event-listener.log
; stdout_logfile_maxbytes=20MB
; stdout_logfile_backups=5
; environment=APP_ENV=production,APP_DEBUG=false

; ==============================================================================
# GROUP DEFINITIONS (For managing related processes together)
; ==============================================================================

[group:erp-queues]
programs=erp-queue-high,erp-queue-default,erp-queue-low,erp-queue-rab-import,erp-queue-ai-classify

[group:erp-horizon-stack]
programs=erp-horizon,erp-scheduler

[group:erp-all]
programs=erp-horizon,erp-scheduler,erp-queue-high,erp-queue-default,erp-queue-low,erp-queue-rab-import,erp-queue-ai-classify

; ==============================================================================
# EVENT LISTENERS (For process lifecycle hooks)
; ==============================================================================

; [eventlistener:erp-queue-restart]
; command=php /var/www/erp-konstruksi/artisan queue:restart
; events=PROCESS_STATE_STOPPED,PROCESS_STATE_FATAL
; ; Restart queue workers on code deploy

; ==============================================================================
# LOG ROTATION (Configure in /etc/logrotate.d/erp-konstruksi)
; ==============================================================================
# /var/log/supervisor/erp-*.log {
#     daily
#     missingok
#     rotate 14
#     compress
#     delaycompress
#     notifempty
#     create 0640 www-data www-data
#     sharedscripts
#     postrotate
#         supervisorctl reload > /dev/null 2>&1 || true
#     endscript
# }

; ==============================================================================
# MONITORING & HEALTH CHECKS
; ==============================================================================
# Supervisor HTTP endpoint: http://127.0.0.1:9001 (user: admin, pass: from env)
# Prometheus exporter available: https://github.com/beancount/supervisor_exporter
# Health check: supervisorctl status | grep -E "RUNNING|FATAL"